# system-wide-scripts/profile.d/vt-env.sh
# /etc/profile.d/vt-env.sh
# VT System Environment Variables
# Version: 1.0
# Author: Echocolate from Sakura-Moon Studio

export vt_sshPort='22'
export vt_wwwUser='www'
export vt_dryRun='0'